# systemd/user/battery-notify.service
[Unit]
Description=Battery Notification Script

[Service]
Type=oneshot
ExecStart=/home/dudel/.config/hypr/local/own_battery_notification.sh

